# artifacts/centrunk.voc.service
[Unit]
Description=Centrunk P25 TRS Service (VOC)
After=syslog.target zerotier-one.service network.target

[Service]
User=root
Type=forking
ExecStart=/opt/dvm/dvmhost -c /opt/dvm/configVOC1.yml
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
